# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:ollama]
command=/usr/local/bin/ollama serve
autostart=true
autorestart=true
stderr_logfile=/var/log/ollama.err.log
stdout_logfile=/var/log/ollama.out.log
environment=OLLAMA_HOST=0.0.0.0:11434
priority=10

[program:postgresql]
command=/usr/local/bin/init-db.sh
user=postgres
autostart=true
autorestart=true
stderr_logfile=/var/log/postgresql.err.log
stdout_logfile=/var/log/postgresql.out.log
priority=20

[program:n8n]
command=/usr/local/bin/n8n
autostart=true
autorestart=true
stderr_logfile=/var/log/n8n.err.log
stdout_logfile=/var/log/n8n.out.log
environment=N8N_HOST=0.0.0.0,N8N_PORT=5678,N8N_PROTOCOL=http,WEBHOOK_URL=http://localhost:5678/,OLLAMA_BASE_URL=http://localhost:11434,DB_TYPE=postgresdb,DB_POSTGRESDB_HOST=localhost,DB_POSTGRESDB_PORT=5432,DB_POSTGRESDB_DATABASE=n8n,DB_POSTGRESDB_USER=n8n,DB_POSTGRESDB_PASSWORD=n8n
priority=30

[program:agent-zero]
command=/exe/run_A0.sh
autostart=true
autorestart=true
stderr_logfile=/var/log/agent-zero.err.log
stdout_logfile=/var/log/agent-zero.out.log
priority=100

[program:searxng]
command=/exe/run_searxng.sh
autostart=true
autorestart=true
stderr_logfile=/var/log/searxng.err.log
stdout_logfile=/var/log/searxng.out.log
priority=40